# supervisord.conf
[supervisord]
nodaemon=true

[program:web]
command=python -u collector.py
directory=/
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:worker]
command=python -u processor.py
directory=/
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

#CMD ["venv/bin/python","-u", "collector.py"]